# root/etc/uci-defaults/30_luci-theme-rubik
#!/bin/sh
uci batch <<-EOF
	set luci.themes.rubik=/luci-static/rubik
	set luci.main.mediaurlbase=/luci-static/rubik
	commit luci
EOF
exit 0
